# nfctd.service.in
[Unit]
Description=Netfilter Conntrack NET-SNMP statistics
After=network.target

[Service]
PIDFile=/var/run/nfctd.pid
ExecStart=@prefix@/bin/nfctd @sysconfdir@/nfctd/nfctd.conf
Restart=on-failure

[Install]
WantedBy=multi-user.target
